# bsc-sniper.service

[Unit]
Description=BSC Sniper Bot 2.0 - Ultimate Money Making Machine
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/bot_files
Environment=PATH=/usr/bin:/bin
Environment=PYTHONPATH=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/bot_files/bsc_sniper_bot_2.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bsc-sniper

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
